# systemd/.config/systemd/user/monitor-hotplug.service
[Unit]
Description=Reload BSPWM

[Service]
Type=oneshot
Environment=DISPLAY=:0
ExecStart= /bin/bash -c "bspc wm -r"
StandardOutput=journal
