# root/etc/udev/rules.d/99-power-energy.rules
SUBSYSTEM=="power_supply", ACTION=="change", RUN+="/home/tiagovla/.local/bin/fan_speed"
